# etc/services.d/xvfb/run
#!/usr/bin/with-contenv bash
set -e
if [[ "$RESOLUTION" == "" ]]; then
    RESOLUTION=1280x1024
fi
Xvfb "$DISPLAY" -screen 0 "$RESOLUTION"
